# git-mirror-auto.timer
[Unit]
Description=Run git-mirror-auto every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
