# git-mirror-auto.timer
[Unit]
Description=Run git-mirror-auto every five minutes

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
